# blended.docker/blended.docker.build/files/profile.d/maven.sh
export M2_HOME=/opt/maven
export PATH=$M2_HOME/bin:$PATH