# blended.docker/blended.docker.build/files/profile.d/scala.sh
export SCALA_HOME=/opt/scala
export PATH=$PATH:$SCALA_HOME/bin
